# assets/systemd/netwatchd.service
[Unit]
Description=NetWatch Statistics Daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/netwatchd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
