# assets/systemd/netwatchd-system.service
[Unit]
Description=NetWatch Statistics Daemon (System-wide)
After=network.target

[Service]
Type=simple
