# assets/systemd/netwatchd-system.service
[Unit]
Description=NetWatch Statistics Daemon (System-wide)
After=network.target

[Service]
Type=simple
User=netwatch
Group=netwatch
ExecStart=/usr/bin/netwatchd --config /etc/netwatch/config.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
